# apps/web-platform/infra/cron-egress-resolve.service
[Unit]
Description=Soleur egress allowlist re-resolve
After=network-online.target
Wants=network-online.target
OnFailure=cron-egress-alarm@%n.service

[Service]
Type=oneshot
# Bounded so a wedged dockerd cannot hold the unit in `activating` forever
# (OnUnitActiveSec timers do not re-fire while activating, and OnFailure=
# never fires for a hang) — a timeout converts the hang into a paged failure.
TimeoutStartSec=120
# Persistent store for the grace-window IP-retention pool (LB-rotation fix).
# systemd creates/owns /var/lib/cron-egress-resolve (root, 0755) BEFORE ExecStart
# and it survives reboots — /run (tmpfs) would wipe the accumulated rotation pool
# and re-open the egress-drop outage for up to one grace window after every boot.
StateDirectory=cron-egress-resolve
# Doppler env wiring — see cron-egress-firewall.service for the rationale.
EnvironmentFile=-/etc/default/inngest-server
Environment=HOME=/root
ExecStart=/bin/sh -c 'D="$(command -v doppler || true)"; if [ -n "$D" ] && [ -n "$DOPPLER_TOKEN" ]; then exec "$D" run --project soleur --config prd -- /usr/local/bin/cron-egress-resolve.sh; else exec /usr/local/bin/cron-egress-resolve.sh; fi'
