# apps/web-platform/infra/cron-egress-resolve.service
[Unit]
Description=Soleur egress allowlist re-resolve (additive-then-prune, #5046 PR-2)
After=network-online.target
Wants=network-online.target
OnFailure=cron-egress-alarm@%n.service

[Service]
Type=oneshot
# doppler (prd) provides SENTRY_* (Crons check-in) + SUPABASE_*/SENTRY_INGEST
# dynamic hosts; the script degrades gracefully when the env is absent.
ExecStart=/usr/bin/doppler run --project soleur --config prd -- /usr/local/bin/cron-egress-resolve.sh
