# apps/web-platform/infra/inngest-cutover-flip.service
[Unit]
Description=Inngest cutover flip oneshot (FSM: armed->flipping->done / rollback; #6178)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Runs as ROOT (no User=): the flip must `systemctl stop/start inngest-server.service`
# and write the /var/lock state slot — a privileged control op, unlike the deploy-scoped
# server/heartbeat units. HOME=/root is set so the Doppler CLI's os.UserHomeDir() init
# does not fail with "$HOME is not defined" (the reason the deploy-user units set User=).
Environment=HOME=/root
# EnvironmentFile carries DOPPLER_TOKEN + DOPPLER_CONFIG_DIR so `doppler run` (and the
# script's `doppler secrets set` flag transitions) authenticate against soleur-inngest/prd.
EnvironmentFile=/etc/default/inngest-server
# doppler run injects INNGEST_CUTOVER_FLIP (the FSM flag) + INNGEST_REDIS_PASSWORD (for
# the loopback redis-cli FLUSHALL/DBSIZE) into the script's env. @@DOPPLER_PROJECT@@ is
# substituted by inngest-bootstrap.sh (dedicated host => soleur-inngest).
ExecStart=/usr/bin/doppler run --project @@DOPPLER_PROJECT@@ --config prd -- /usr/local/bin/inngest-cutover-flip.sh

[Install]
WantedBy=multi-user.target
