# apps/web-platform/infra/orphan-reaper.service
[Unit]
Description=Orphaned workspace directory reaper

[Service]
Type=oneshot
ExecStart=/usr/local/bin/orphan-reaper.sh
