# apps/web-platform/infra/orphan-reaper.timer
[Unit]
Description=Run orphan reaper every 6 hours

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
