# apps/web-platform/infra/webhook.service
[Unit]
Description=Webhook deploy listener
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/default/webhook-deploy
# -ip 0.0.0.0 (was 127.0.0.1): the private-net peer must reach /hooks/deploy-peer
# for the 2-host deploy fan-out (#5274 Phase 3, ADR-068). SAFE because
# hcloud_firewall.web default-denies 9000 on the PUBLIC interface (only 22/80/443/
# icmp are allowed) — so only loopback + the private net (10.0.1.0/24) can connect.
# That default-deny is load-bearing; firewall-9000-deny.test.sh pins it.
ExecStart=/usr/local/bin/webhook -verbose -hooks /etc/webhook/hooks.json -port 9000 -ip 0.0.0.0
Restart=on-failure
RestartSec=5
User=deploy
Group=deploy
# NoNewPrivileges omitted: ci-deploy.sh requires sudo for chown (#967)
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
# Bootstrap path needs to install systemd units + binaries + env files under
# `sudo` (inngest-bootstrap.sh + any future `deploy <component>` handler).
# ProtectSystem=strict makes /etc/ /usr/ read-only at the mount-namespace
# level — sudo elevates UID but does NOT escape the namespace. Without these
# paths in ReadWritePaths, `cat > /etc/systemd/system/...` and
# `chown root:deploy /etc/default/...` fail with `Read-only file system` or
# silently mangled error messages (`invalid group: 'root:deploy'` because the
# initial `printf >` redirect actually failed but chown still ran).
# Surfaced 2026-05-19 via #4017 — PR-F's inngest substrate was never installed
# in production because the webhook deploy of the bootstrap image failed at
# the chown step inside this sandbox. The /var/lib/inngest entry is required
# because mkdir+chown for SQLite's data dir runs ahead of the systemd unit's
# own ReadWritePaths.
# /etc/vector + /var/lib/vector added 2026-05-21 (TR9 PR-5): inngest-bootstrap.sh
# now also installs Vector (observability shipper) alongside the inngest CLI;
# its config + state dirs need bootstrap-time write access from inside this
# namespace. Surfaced via v1.1.1 deploy diagnostic capture
# (reason="mkdir: cannot create directory /etc/vector: Read-only file system").
# `-` prefix on /var/lib/vector + /etc/vector marks them optional:
# systemd silently ignores the path if absent, instead of refusing to set up
# the mount namespace (which on prior PR #4257 took webhook.service down
# entirely until the dirs were created). Once inngest-bootstrap.sh creates
# them on first v1.1.x deploy, they become real ReadWritePaths.
ReadWritePaths=/mnt/data /var/lock /etc/systemd/system /etc/default /etc/sudoers.d /etc/webhook /var/lib/inngest -/var/lib/vector -/etc/vector /usr/local/bin
ReadOnlyPaths=/etc/default/webhook-deploy
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target
